# deploy/systemd/reconciliation-api.service
[Unit]
Description=ERP / 1C Reconciliation API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/recon_erp_1c
ExecStart=/opt/recon_erp_1c/.venv/bin/python -m recon_erp_1c.interfaces.http.api
Restart=always
RestartSec=2
EnvironmentFile=-/etc/recon_erp_1c.env

[Install]
WantedBy=multi-user.target
